# autoresearch.conf
[program:autoresearch]
directory=/home/ubuntu/autoresearch
command=/home/ubuntu/.local/bin/uv run daemon.py
autostart=true
autorestart=true
stderr_logfile=/home/ubuntu/autoresearch/daemon.err.log
stdout_logfile=/home/ubuntu/autoresearch/daemon.out.log
user=ubuntu
environment=GOOGLE_API_KEY="YOUR_API_KEY_HERE"
